# imx/meta-fsl-arm/recipes-connectivity/bluez5/bluez5/bluetooth.in
#!/bin/sh
#
# start/stop bluetooth daemon.

### BEGIN INIT INFO
# Provides:          bluetooth
# Required-Start:    $network $udev
# Required-Stop:     $network
# Default-Start:     S 2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Bluetooth daemon
# Description:       starts bluetooth daemon
### END INIT INFO

DAEMON=/usr/lib/bluez5/bluetooth/bluetoothd
PIDFILE=/var/run/bluetooth.pid
DESC="Linux Bluetooth daemon"

if [ -f /etc/default/bluetooth ] ; then
        . /etc/default/bluetooth
fi

set -e

do_start() {
        start-stop-daemon --start -b -x $DAEMON -p $PIDFILE
	/usr/bin/udevadm settle --exit-if-exists=$PIDFILE
	/usr/bin/udevadm trigger --action=add --subsystem-match=rfkill
}

do_stop() {
        start-stop-daemon --stop --name bluetoothd --quiet -p $PIDFILE
}

case "$1" in
  start)
        echo "Starting $DESC"
        do_start
        ;;
  stop)
        echo "Stopping $DESC"
        do_stop
        ;;
  restart|force-reload)
        echo "Restarting $DESC"
        do_stop
        sleep 1
        do_start
        ;;
  *)
        echo "Usage: $0 {start|stop|restart|force-reload}" >&2
        exit 1
        ;;
esac

exit 0

