# imx/meta-bsp/recipes-core/initscripts/arm/mx6/gpuload.service
[Unit]
Description=Formerly rc_gpu.S

[Service]
ExecStart=/usr/bin/gpuload

[Install]
WantedBy=multi-user.target 
