# imx/meta-bsp/recipes-core/systemd/systemd-gpuconfig/gpuconfig.service
[Unit]
Description=Configuration for i.MX GPU (Former rc_gpu.S)

[Service]
ExecStart=/etc/gpuconfig

[Install]
WantedBy=multi-user.target
