# imx/meta-bsp/recipes-core/initscripts/arm/aarch64/alignment.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides: alignment
# Required-Start:    mountkernfs
# Required-Stop:     mountkernfs
# Default-Start:     S
# Default-Stop:
### END INIT INFO
