# meta-bsp/recipes-bsp/firmware-imx/files/hdmi-firmware.service
[Unit]
Description=User Space HDMI Firmware Loading

[Service]
ExecStart=/etc/hdmi

[Install]
WantedBy=multi-user.target
