# meta-bsp/recipes-bsp/firmware-imx/files/regulatory-firmware.service
[Unit]
Description=User Space Regulatory Firmware Loading

[Service]
ExecStart=/etc/regulatory

[Install]
WantedBy=multi-user.target
