# meta-bsp/recipes-bsp/firmware-imx/files/sdma-firmware.service
[Unit]
Description=User Space SDMA Firmware Loading

[Service]
ExecStart=/etc/sdma

[Install]
WantedBy=multi-user.target
