# imx/meta-bsp/recipes-core/systemd/systemd/userspace-firmware.service
[Unit]
Description=User Space Firmware Loading

[Service]
ExecStart=/etc/sdma

[Install]
WantedBy=multi-user.target
