# imx/meta-bsp/recipes-bsp/firmware-imx/files/epdc-firmware.service
[Unit]
Description=User Space Epdc Firmware Loading

[Service]
ExecStart=/etc/epdc

[Install]
WantedBy=multi-user.target
