# imx/meta-bsp/recipes-bsp/firmware-imx/files/hdmitx-firmware.service
[Unit]
Description=User Space HDMI TX Firmware Loading

[Service]
ExecStart=/etc/hdmitx

[Install]
WantedBy=multi-user.target
